# meta-rockchip-rv1106/recipes-graphics/luckfox-lvgl-demo/files/luckfox-lvgl-demo.service
[Unit]
Description=Luckfox Pico LVGL framebuffer demo
After=multi-user.target

[Service]
Type=simple
ExecStart=/usr/bin/luckfox-lvgl-demo
Restart=on-failure
RestartSec=2

[Install]
WantedBy=multi-user.target
